# attestation-verifier/src/build/linux/tagent/tagent.service
[Unit]
Description=TrustAgent Service
Requires=tagent_init.service
After=tagent_init.service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/trustagent/bin
ExecStart=/opt/trustagent/bin/tagent startService
ExecStop=/bin/kill -s HUP $MAINPID
TimeoutStartSec=45
Restart=on-failure
PermissionsStartOnly=true
User=tagent

[Install]
WantedBy=multi-user.target
