# attestation-verifier/src/build/linux/tagent/tagent_init.service
[Unit]
Description=TrustAgent Initialization
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/trustagent/bin
ExecStart=/opt/trustagent/bin/tagent init

[Install]
WantedBy=multi-user.target
