# attestation-verifier/utils/tools/bkc/build/installer/bkc-reboot.service
[Unit]
Description=bkc-tool reboot

[Service]
Type=idle
WorkingDirectory=/opt/bkc-tool
ExecStart=/opt/bkc-tool/bin/reboot.sh

[Install]
WantedBy=multi-user.target
